# net-fs/glusterfs/files/glusterfsd.initd
#!/sbin/runscript
# Copyright 1999-2008 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: $

start() {
	ebegin "Starting GlusterFSd"
	/usr/sbin/glusterfsd --pidfile=/var/run/glusterfsd.pid ${GLUSTERFSD_OPTS}
	eend $?
}

stop() {
	ebegin "Stopping GlusterFSd"
	start-stop-daemon --stop --quiet --pidfile /var/run/glusterfsd.pid
	eend $?
}
